# rootfs/etc/s6-overlay/s6-rc.d/init-keepalived/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Keepalived
# Ensures keepalived.conf is copied into the add one
# ==============================================================================

# Migrate add-on data from the Home Assistant config folder,
# to the add-on configuration folder.
if ! bashio::fs.file_exists '/config/keepalived.conf' \
    && bashio::fs.file_exists '/homeassistant/keepalived.conf'; then
    cp /homeassistant/keepalived.conf /config/keepalived.conf \
        || bashio::exit.nok "Failed to copy Keepalived configuration"
fi
